# greeter/lightdm-kde-greeter-wifikeeper.service.in
# SPDX-FileCopyrightText: None
# SPDX-License-Identifier: CC0-1.0

[Unit]
Description=LightDM KDE Greeter, to create a wifi connection that will last a little longer than greeter

[Service]
Type=simple
ExecStart=@KDE_INSTALL_FULL_LIBEXECDIR@/lightdm-kde-greeter-wifikeeper
